# xdr-trace-client_prod_1.8.0-bullseye_amd64/etc/cron.d/trace
# Copyright : Whizhack Technologies, India
# Authors : Team Devops

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Scheduled software restarts and system cleanups.
33 0,6,12,18 * * *     root    /opt/zh-xdr-trace-client/management/reset

# Scheduled software updates.
33 3 * * 1              root    /opt/zh-xdr-trace-client/management/update

# Scheduled health log deletion.
5 4 */14 * *            root    rm /var/log/trace/health.log; touch /var/log/trace/health.log;

# Scheduled credential updates after each reboot.
@reboot                 root    /usr/local/bin/trace credentialRefresh;

