# installers/debian/xdr-hids-manager-client_dev_1.8.0-bullseye_amd64/lib/systemd/system/xdr-hids-health.service
[Unit]
Description=xdr-hids-health
 
[Service]
Restart=always
RestartSec=10
TimeoutSec=infinity

ExecStart=/usr/local/bin/healthcheck 
 
[Install]
WantedBy=multi-user.target

